# examples/oterm/oterm.socket
[Unit]
Description=Onion web terminal emulator.

[Socket]
ListenStream=8080

[Install]
WantedBy=oterm.target
