# examples/oterm/oterm.service
[Service]
ExecStart=/usr/local/bin/oterm

